# crontab
*/5 * * * * cd /workspaces/data_process && /usr/local/bin/python manage.py delete_old_files